# deploy/anilax-api.service
[Unit]
Description=Anilax Software API (Express + PostgreSQL)
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/anilax-software/anilax-software-backend
EnvironmentFile=/var/www/anilax-software/anilax-software-backend/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
